# etc/cron.daily/install-exim-certificates
#!/bin/bash

set -eu

/etc/bin/install-lecerts irc.darwin.network Debian-exim nogroup /etc/eximkeys/tls

# TODO is this necessary?
systemctl restart exim4.service

exit 0
