# etc/cron.daily/install-stunnel-certificates
#!/bin/bash

set -eu

/etc/bin/install-lecerts irc.darwin.network oragono nogroup /home/oragono/tls

systemctl reload proxytunnel.service

exit 0
