# etc/cron.daily/install-stunnel-certificates
#!/bin/bash

set -eu

export LE_SOURCEDIR=/etc/certificates
/etc/bin/install-lecerts irc.darwin.network stunnel4 stunnel4 /etc/stunnel_certificates

systemctl reload proxytunnel.service

exit 0
